# salt/roots/salt/templates/core.ini
[program:core]
command = /var/themis/quals/core/bin/app
directory = /var/themis/quals/core
process_name = %(program_name)s-%(process_num)d
numprocs = {{ processes }}
autostart = true
autorestart = true
user = vagrant
redirect_stderr = true
environment = SESSION_SECRET="{{ secret }}",CORS_ORIGIN="{{ origin }}",SERVER_PORT="30%(process_num)02d"
