# salt/roots/salt/templates/queue.ini
[program:queue]
command = /var/themis/quals/core/bin/queue
directory = /var/themis/quals/core
process_name = %(program_name)s-%(process_num)d
numprocs = {{ processes }}
autostart = true
autorestart = true
user = {{ user }}
redirect_stderr = true
environment = MONGODB_URI="{{ mongodb_uri }}",LOGOS_DIR="{{ logos_dir }}",DOMAIN="{{ domain }}",MANDRILL_API_KEY="{{ mandrill_api_key }}",EMAIL_SENDER="{{ email_sender }}"
